# supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:homeassistant]
command=hass -c /data/homeassistant
user=solar
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/homeassistant.log
stderr_logfile=/var/log/supervisor/homeassistant.err.log

[program:influxdb]
command=influxd --config /etc/influxdb/config.yaml
user=solar
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/influxdb.log
stderr_logfile=/var/log/supervisor/influxdb.err.log

[program:grafana]
command=/usr/share/grafana/bin/grafana-server --config=/etc/grafana/grafana.ini --homepath=/usr/share/grafana
user=solar
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/grafana.log
stderr_logfile=/var/log/supervisor/grafana.err.log

[program:mosquitto]
command=mosquitto -c /etc/mosquitto/mosquitto.conf
user=solar
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/mosquitto.log
stderr_logfile=/var/log/supervisor/mosquitto.err.log

[program:node-red]
command=node-red -u /data/node-red --settings /etc/node-red/settings.js
user=solar
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/node-red.log
stderr_logfile=/var/log/supervisor/node-red.err.log

[program:uptime-kuma]
command=node /usr/share/uptime-kuma/server/server.js
directory=/usr/share/uptime-kuma
user=solar
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/uptime-kuma.log
stderr_logfile=/var/log/supervisor/uptime-kuma.err.log
environment=PORT="3001"

[program:open-meteo]
command=/usr/local/bin/open-meteo-wrapper.sh
user=solar
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/open-meteo.log
stderr_logfile=/var/log/supervisor/open-meteo.err.log

[program:energy-guard]
command=/usr/local/bin/energy-guard --config /etc/energy-guard/config.json
user=solar
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/energy-guard.log
stderr_logfile=/var/log/supervisor/energy-guard.err.log
